# service/nginx/run
#!/bin/sh

exec chpst -u nginx /usr/sbin/nginx -c $PWD/nginx.conf
